# systemd/user/tmux.service.d/override.conf
# using `sysemctl --user edit tmux.service` to edit config file
[Service]
Environment="PATH=/home/linuxbrew/.linuxbrew/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
