# systemd/user/tmux.service
[Unit]
Description=tmux default session (detached)
Documentation=man:tmux(1)

[Service]
Type=forking
Environment=DISPLAY=:0
ExecStart=/home/linuxbrew/.linuxbrew/bin/tmux new-session -d

ExecStop=/home/soros/.tmux/plugins/tmux-resurrect/scripts/save.sh
ExecStop=/home/linuxbrew/.linuxbrew/bin/tmux kill-server
KillMode=control-group

RestartSec=2

[Install]
WantedBy=default.target
