# systemd/system/pkg-ContainerManager-dockerd.service.d/http-proxy.conf
[Service]
Environment="HTTP_PROXY=http://nuc.sorosliu.xyz:7890/"
Environment="HTTPS_PROXY=http://nuc.sorosliu.xyz:7890/"
